# build/packaging/suse-uptime-tracker.service
[Unit]
Description=Run SUSE uptime tracker
Wants=suse-uptime-tracker.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/suse-uptime-tracker
