# build/packaging/suse-uptime-tracker.timer
[Unit]
Description=Schedule uptime tracking every 15 minutes

[Timer]
# Run this timer every 15 minutes
OnCalendar=*:0/15
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
